# rspi/etc/rc.local
#!/bin/sh -e

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

# Startup Redirect HTTP --> HTTPS app
/usr/local/bin/node /home/pi/node/app/redirect/app.js 2>&1 & 

# Startup dropkey app
su -l pi -c '/home/pi/node/dropkey.git/hooks/post-receive 2>&1 &'

exit 0


